# configs/lora_gateway.service
[Unit]
After=lora.service

[Service]
ExecStart=/usr/local/bin/lora_gateway.sh

[Install]
WantedBy=default.target
